# source/usr/local/emhttp/plugins/dwsnap/misc/dwsnap-maintenance.cron
# SnapRAID Maintenance:
* * * * * /usr/bin/snapraid-cron 1>/dev/null 2>(logger -t "snapraid-cron")
