# source/usr/local/emhttp/plugins/dwsnap/misc/dwsnap-maintenance.cron
# SnapRAID Maintenance:
* * * * * /usr/bin/snapraid-cron > /dev/null 2>&1
